# openobserve_log_shipper/rootfs/etc/services.d/openobserve/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e

OPENOBSERVE_URL="$(bashio::config 'openobserve_url')"
OPENOBSERVE_USERNAME="$(bashio::config 'openobserve_username')"
OPENOBSERVE_PASSWORD="$(bashio::config 'openobserve_password')"
LOG_PATH="$(bashio::config 'log_path')"

if [[ -z "${OPENOBSERVE_URL}" ]]; then
  bashio::exit.nok "openobserve_url must be set (e.g. http://host:5080/api/default/home_assistant_supervisor/_multi)"
fi

AUTH_HEADER=""
if [[ -n "${OPENOBSERVE_USERNAME}" && -n "${OPENOBSERVE_PASSWORD}" ]]; then
  AUTH_HEADER="Basic $(printf '%s:%s' "${OPENOBSERVE_USERNAME}" "${OPENOBSERVE_PASSWORD}" | base64 | tr -d '\n')"
fi

# Parse host, port, uri from full ingest URL (http://host:port/path)
if [[ -n "${OPENOBSERVE_URL}" ]]; then
  O2_SCHEME="${OPENOBSERVE_URL%%://*}"
  REST="${OPENOBSERVE_URL#*://}"
  O2_HOST="${REST%%/*}"
  O2_URI="/${REST#*/}"
  if [[ "${O2_HOST}" == *:* ]]; then
    O2_PORT="${O2_HOST##*:}"
    O2_HOST="${O2_HOST%%:*}"
  elif [[ "${O2_SCHEME}" == "https" ]]; then
    O2_PORT="443"
  else
    O2_PORT="80"
  fi
  TLS="Off"
  if [[ "${O2_SCHEME}" == "https" ]]; then
    TLS="On"
  fi
fi

export O2_HOST O2_PORT O2_URI TLS LOG_PATH AUTH_HEADER

envsubst '${O2_HOST} ${O2_PORT} ${O2_URI} ${TLS} ${LOG_PATH} ${AUTH_HEADER}' \
  < /etc/fluent-bit/fluent-bit.conf.template \
  > /etc/fluent-bit/fluent-bit.conf

bashio::log.info "Starting Fluent Bit → ${OPENOBSERVE_URL}"
exec fluent-bit -c /etc/fluent-bit/fluent-bit.conf
